# anchorbot.service
[Unit]
Description=Foo

[Service]
ExecStart=/usr/bin/python3 /root/anchorChecker